# khal_notifier/khal_notifier.service
[Unit]
Description=Sync calendars and show upcoming events

[Service]
Type=oneshot
ExecStart=/opt/tools/khal_notifier
